# templates/default/godns.init.conf.erb
# Generated by Chef
#

description "GoDNS - DNS server backed by redis"
author "leo@unbekandt.eu"

start on runlevel [2345]
stop on runlevel [016]

respawn

exec <%= "#{@target} -c #{node['godns']['config_path']} %>

